# systemd/claimvm.slice
[Unit]
Description=ClaimVM load slice

[Slice]
CPUWeight=idle
Nice=19
